# ansible/roles/quarkslab.systemd/files/irma.target
[Unit]
Description=IRMA target
Requires=multi-user.target
After=multi-user.target
